# arjunbishnoi.service
[Unit]
Description=Arjun Bishnoi Personal Website
Documentation=https://github.com/yourusername/arjunbishnoi.com
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/arjunbishnoi.com
ExecStart=/usr/bin/node .output/server/index.mjs
Restart=on-failure
RestartSec=10
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target 